# server/crontab
# crontab -e

@reboot /home/pi/lilygo-t5/server/server.py 2>&1 >> /home/pi/lilygo-t5-server.log
